# services/codex-usage-collector/codex-usage-collector.service.example
[Unit]
Description=ZXLab Codex Usage Collector
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/zxlab-codex-usage
EnvironmentFile=%h/.config/zxlab-codex-usage.env
Environment=PATH=%h/.local/node-v22/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env node %h/zxlab-codex-usage/server.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
